# vncserver/vncserver.socket
[Unit]
Description=XVNC Server on port 5901

[Socket]
ListenStream=[::]:5901

[Install]
WantedBy=sockets.target
