# docker-registry-ui/docker-registry-ui.socket
[Unit]
Description=docker-registry

[Socket]
ListenStream=0.0.0.0:3001

[Install]
WantedBy=sockets.target
