# calibre-web/calibre-web.socket
[Unit]
Description=grafana

[Socket]
ListenStream=5000

[Install]
WantedBy=sockets.target
