# kopia/kopia.socket
[Unit]
Description=kopia

[Socket]
ListenStream=55515

[Install]
WantedBy=sockets.target
