# gitserver/gitserver.socket
[Unit]
Description=grafana

[Socket]
ListenStream=3558

[Install]
WantedBy=sockets.target
