# photoprism/photoprism.socket
[Unit]
Description=grafana

[Socket]
ListenStream=2342

[Install]
WantedBy=sockets.target
